# es-cron/es-indexdelete.cron
30 1 * * * root /opt/isys/yz-es/es-indexdelete.sh >> /data/logs/es-indexdelete-main.log 2>&1
